# bash/.bash/profile.d/05-thefuck
# bash profile: thefuck
#

if ! hash thefuck > /dev/null 2>&1; then
    return
fi

eval "$(thefuck --alias)"


# vim: set syntax=sh filetype=sh ts=4 sw=4 tw=0 et :
