# bash/.bash/profile.d/05-tmuxinator
# bash profile: tmuxinator
#

# Bash completions
source /usr/lib/ruby/gems/2.5.0/gems/tmuxinator-0.11.1/completion/tmuxinator.bash


# vim: set syntax=sh filetype=sh ts=4 sw=4 tw=0 et :
