# bash/.bash/profile.d/05-thefuck
# bash profile: thefuck
#

eval "$(thefuck --alias)"


# vim: set syntax=sh filetype=sh ts=4 sw=4 tw=0 et :
